# systemd/grive-timer@.service.in
[Unit]
Description=Google drive sync (executed by timer unit)
After=network-online.target

[Service]
ExecStart=@GRIVE_SYNC_SH_BINARY@ sync "%i"
